# packaging/LaunchAgents/com.richardhuo.aiquotawatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.richardhuo.aiquotawatch</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Applications/AI 额度观察.app/Contents/MacOS/AIQuotaWatch</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>WorkingDirectory</key>
    <string>__HOME__</string>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/AIQuotaWatch.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/AIQuotaWatch.err.log</string>
</dict>
</plist>
